# Kubernetes_project/srcs/phpmyadmin/srcs/supervisord.conf
# -S option : run with built-in web server. -S <addr>:<port>
# Specify the root of the document for the internal web server
# telegraf &
# php -S 0.0.0.0:5000 -t /var/www/phpmyadmin/
# nginx -g "daemon off;"

[supervisord]
nodaemon=true

[program:php-fpm7]
command=php-fpm7 --nodaemonize
autostart=true
autorestart=true

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true

[program:telegraf]
command=telegraf
autostart=true
autorestart=true