# eft_services/srcs/influxdb/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:influxdb]
command=influxd
user=root
autostart=true
autorestart=true

[program:telegraf]
command=telegraf
user=root
startsecs=120
autostart=true
autorestart=true